# cerberus
#!/usr/bin/env bash

### BEGIN INIT INFO
# Provides:          cerberus
# Required-Start:
# Required-Stop:
# Should-Start:
# Default-Start:
# Default-Stop:
# Short-Description: Transparent Proxy through TOR.
### END INIT INFO

# Cerberus v1.0.0AlfaStable - Transparent Proxy through TOR
# Engineered by 0xb0rn3 | 0xbv1
# Repository: https://github.com/0xb0rn3/cerberus
# 
# Based on AnonSurf, adapted for Arch-based systems
# License: GNU General Public License v3

export BLUE='\033[1;94m'
export GREEN='\033[1;92m'
export RED='\033[1;91m'
export RESETCOLOR='\033[1;00m'

# Destinations you don't want routed through Tor
TOR_EXCLUDE="192.168.0.0/16 172.16.0.0/12 10.0.0.0/8"

# The UID Tor runs as (Arch uses 'tor' user)
TOR_UID="tor"

# Tor's TransPort
TOR_PORT="9040"

# Check if systemd-resolved is active
if systemctl is-active --quiet systemd-resolved; then
    resolved_support=true
else
    resolved_support=false
fi

function init {
    echo -e -n " $GREEN*$BLUE killing dangerous applications$RESETCOLOR\n"
    killall -q chrome dropbox firefox chromium skype thunderbird transmission deluge pidgin telegram discord slack signal-desktop
    
    echo -e -n " $GREEN*$BLUE cleaning cache elements$RESETCOLOR\n"
    # Clean browser caches if bleachbit is available
    if command -v bleachbit > /dev/null; then
        bleachbit -c firefox.cache chromium.cache google_chrome.cache &> /dev/null
    fi
}

function starti2p {
    echo -e -n " $GREEN*$BLUE starting I2P services$RESETCOLOR\n"
    systemctl stop tor
    
    # Modify DNS settings
    if [ "$resolved_support" = true ]; then
        # Backup current resolved config
        cp /etc/systemd/resolved.conf /etc/systemd/resolved.conf.bak
        # Configure DNS
        echo -e "[Resolve]\nDNS=127.0.0.1\nFallbackDNS=209.222.18.222 209.222.18.218" > /etc/systemd/resolved.conf
        systemctl restart systemd-resolved
    else
        cp /etc/resolv.conf /etc/resolv.conf.bak
        echo -e 'nameserver 127.0.0.1\nnameserver 209.222.18.222\nnameserver 209.222.18.218' > /etc/resolv.conf
    fi
    
    # Start I2P if installed
    if command -v i2prouter > /dev/null; then
        sudo -u i2p i2prouter start
        sleep 2
        if command -v xdg-open > /dev/null; then
            xdg-open 'http://127.0.0.1:7657/home'
        fi
    else
        echo -e " $RED*$BLUE I2P not installed$RESETCOLOR\n"
    fi
}

function stopi2p {
    echo -e -n " $GREEN*$BLUE stopping I2P services$RESETCOLOR\n"
    if command -v i2prouter > /dev/null; then
        sudo -u i2p i2prouter stop
    fi
    
    # Restore DNS settings
    if [ "$resolved_support" = true ]; then
        if [ -e /etc/systemd/resolved.conf.bak ]; then
            mv /etc/systemd/resolved.conf.bak /etc/systemd/resolved.conf
            systemctl restart systemd-resolved
        fi
    else
        if [ -e /etc/resolv.conf.bak ]; then
            mv /etc/resolv.conf.bak /etc/resolv.conf
        fi
    fi
}

function disable_ipv6() {
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Disabling IPv6 services:$RESETCOLOR\n"
    
    # Remove current IPv6 settings if it exists
    if [ -f /etc/sysctl.d/98-cerberus.conf ]; then
        rm /etc/sysctl.d/98-cerberus.conf
    fi
    
    # Add comprehensive IPv6 disable configuration
    cat << EOF >> /etc/sysctl.d/98-cerberus.conf
# Disable IPv6 - Added by Cerberus
net.ipv6.conf.all.disable_ipv6 = 1
net.ipv6.conf.default.disable_ipv6 = 1
net.ipv6.conf.lo.disable_ipv6 = 1
EOF
    
    # Apply settings
    if ! sysctl -p /etc/sysctl.d/98-cerberus.conf > /dev/null; then
        echo -e "$RED Error applying sysctl settings$RESETCOLOR"
        return 1
    fi
    
    echo -e "$GREEN IPv6 disabled$RESETCOLOR"
    return 0
}

function enable_ipv6() {
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Enabling IPv6 services:$RESETCOLOR\n"
    
    # Remove the config file if it exists
    if [ -f /etc/sysctl.d/98-cerberus.conf ]; then
        rm -f /etc/sysctl.d/98-cerberus.conf
    fi
    
    # Reload all sysctl settings
    if ! sysctl --system > /dev/null; then
        echo -e "$RED Error reloading sysctl settings$RESETCOLOR"
        return 1
    fi
    
    # Restart network services (Arch specific)
    systemctl reload NetworkManager > /dev/null 2>&1
    systemctl restart systemd-resolved > /dev/null 2>&1
    
    echo -e "$GREEN IPv6 enabled$RESETCOLOR"
    return 0
}

function ip {
    echo -e "\nCurrent IP address:\n"
    sleep 1
    if command -v curl > /dev/null; then
        curl -s "https://api.ipify.org"
    elif command -v wget > /dev/null; then
        wget -qO- "https://api.ipify.org"
    else
        echo -e "$RED No curl or wget available$RESETCOLOR"
    fi
    echo -e "\n\n----------------------------------------------------------------------"
}

function start {
    # Make sure only root can run this script
    if [ $(id -u) -ne 0 ]; then
        echo -e -e "\n$GREEN[$RED!$GREEN] $RED This script must be run as root$RESETCOLOR\n" >&2
        exit 1
    fi
    
    # Check if Tor is installed
    if ! command -v tor > /dev/null; then
        echo -e "\n$GREEN[$RED!$GREEN]$RED Tor is not installed. Install with: pacman -S tor$RESETCOLOR\n" >&2
        exit 1
    fi
    
    disable_ipv6
    
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Starting Cerberus anonymous mode:$RESETCOLOR\n"
    
    # Start Tor if not running
    if ! systemctl is-active --quiet tor; then
        echo -e " $RED*$BLUE Tor is not running! Starting it$RESETCOLOR\n"
        systemctl reload NetworkManager > /dev/null 2>&1
        systemctl start tor
        sleep 2
    fi
    
    # Backup current iptables rules
    if ! [ -f /etc/iptables/iptables.rules.bak ]; then
        mkdir -p /etc/iptables
        iptables-save > /etc/iptables/iptables.rules.bak
        echo -e " $GREEN*$BLUE Saved iptables rules$RESETCOLOR\n"
    fi
    
    # Flush existing rules
    iptables -F
    iptables -t nat -F
    
    # Configure DNS
    if [ "$resolved_support" = true ]; then
        cp /etc/systemd/resolved.conf /etc/systemd/resolved.conf.bak
        echo -e "[Resolve]\nDNS=127.0.0.1\nFallbackDNS=209.222.18.222 209.222.18.218" > /etc/systemd/resolved.conf
        systemctl restart systemd-resolved
    else
        cp /etc/resolv.conf /etc/resolv.conf.bak
        echo -e 'nameserver 127.0.0.1\nnameserver 209.222.18.222\nnameserver 209.222.18.218' > /etc/resolv.conf
    fi
    
    # Set iptables NAT rules
    iptables -t nat -A OUTPUT -m owner --uid-owner $TOR_UID -j RETURN
    iptables -t nat -A OUTPUT -p udp --dport 53 -j REDIRECT --to-ports 53
    iptables -t nat -A OUTPUT -p tcp --dport 53 -j REDIRECT --to-ports 53
    iptables -t nat -A OUTPUT -p udp -m owner --uid-owner $TOR_UID -m udp --dport 53 -j REDIRECT --to-ports 53
    
    # Resolve .onion domains mapping 10.192.0.0/10 address space
    iptables -t nat -A OUTPUT -p tcp -d 10.192.0.0/10 -j REDIRECT --to-ports 9040
    iptables -t nat -A OUTPUT -p udp -d 10.192.0.0/10 -j REDIRECT --to-ports 9040
    
    # Exclude local addresses
    for NET in $TOR_EXCLUDE 127.0.0.0/9 127.128.0.0/10; do
        iptables -t nat -A OUTPUT -d $NET -j RETURN
    done
    
    # Redirect all other output through TOR
    iptables -t nat -A OUTPUT -p tcp --syn -j REDIRECT --to-ports $TOR_PORT
    iptables -t nat -A OUTPUT -p udp -j REDIRECT --to-ports $TOR_PORT
    iptables -t nat -A OUTPUT -p icmp -j REDIRECT --to-ports $TOR_PORT
    
    # Accept already established connections
    iptables -A OUTPUT -m state --state ESTABLISHED,RELATED -j ACCEPT
    
    # Exclude local addresses
    for NET in $TOR_EXCLUDE 127.0.0.0/8; do
        iptables -A OUTPUT -d $NET -j ACCEPT
    done
    
    # Allow only tor output
    iptables -A OUTPUT -m owner --uid-owner $TOR_UID -j ACCEPT
    iptables -A OUTPUT -j REJECT
    
    echo -e "$GREEN *$BLUE All traffic redirected through Tor$RESETCOLOR\n"
    echo -e "$GREEN[$BLUE i$GREEN ]$BLUE Cerberus tunnel active$RESETCOLOR\n"
}

function stop {
    # Make sure only root can run this script
    if [ $(id -u) -ne 0 ]; then
        echo -e "\n$GREEN[$RED!$GREEN] $RED This script must be run as root$RESETCOLOR\n" >&2
        exit 1
    fi
    
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Stopping Cerberus anonymous mode:$RESETCOLOR\n"
    
    # Flush rules
    iptables -F
    iptables -t nat -F
    echo -e " $GREEN*$BLUE Cleared iptables rules$RESETCOLOR\n"
    
    # Restore iptables rules
    if [ -f /etc/iptables/iptables.rules.bak ]; then
        iptables-restore < /etc/iptables/iptables.rules.bak
        rm /etc/iptables/iptables.rules.bak
        echo -e " $GREEN*$BLUE Restored iptables rules$RESETCOLOR"
    fi
    
    # Restore DNS settings
    if [ "$resolved_support" = true ]; then
        if [ -e /etc/systemd/resolved.conf.bak ]; then
            mv /etc/systemd/resolved.conf.bak /etc/systemd/resolved.conf
            systemctl restart systemd-resolved
        fi
    else
        if [ -e /etc/resolv.conf.bak ]; then
            mv /etc/resolv.conf.bak /etc/resolv.conf
        fi
    fi
    
    # Stop Tor
    systemctl stop tor
    
    # Re-enable IPv6
    enable_ipv6
    
    echo -e " $GREEN*$BLUE Cerberus stopped$RESETCOLOR\n"
}

function change {
    systemctl reload tor
    sleep 2
    echo -e " $GREEN*$BLUE Tor reloaded - identity changed$RESETCOLOR\n"
}

function status {
    if systemctl is-active --quiet tor; then
        echo -e "$GREEN Tor is running$RESETCOLOR"
        
        # Check if iptables rules are active
        if iptables -t nat -L OUTPUT | grep -q "REDIRECT"; then
            echo -e "$GREEN Cerberus tunnel is active$RESETCOLOR"
        else
            echo -e "$RED Cerberus tunnel is not active$RESETCOLOR"
        fi
    else
        echo -e "$RED Tor is not running$RESETCOLOR"
    fi
}

case "$1" in
    start)
        init
        start
    ;;
    stop)
        init
        stop
    ;;
    change)
        change
    ;;
    status)
        status
    ;;
    myip|ip)
        ip
    ;;
    starti2p)
        starti2p
    ;;
    stopi2p)
        stopi2p
    ;;
    restart)
        $0 stop
        sleep 1
        $0 start
    ;;
    *)
        echo -e "
$RED╔═══════════════════════════════════════════════════════════════════════════════╗
$RED║                    $GREEN Cerberus v1.0.0AlfaStable                                 $RED║
$RED║                    $BLUE Engineered by 0xb0rn3 | 0xbv1                            $RED║
$RED║                    $BLUE https://github.com/0xb0rn3/cerberus                      $RED║
$RED╚═══════════════════════════════════════════════════════════════════════════════╝$RESETCOLOR

Usage: cerberus {start|stop|restart|change|status|myip|starti2p|stopi2p}

$GREEN start$BLUE     - Start system-wide anonymous tunneling through TOR
$GREEN stop$BLUE      - Stop anonymous tunneling and restore original settings  
$GREEN restart$BLUE   - Restart Cerberus (stop + start)
$GREEN change$BLUE    - Change TOR identity (get new exit node)
$GREEN status$BLUE    - Check if Cerberus is running
$GREEN myip$BLUE      - Show current IP address
$GREEN starti2p$BLUE  - Start I2P services
$GREEN stopi2p$BLUE   - Stop I2P services
$RESETCOLOR" >&2
        exit 1
    ;;
esac

echo -e $RESETCOLOR
exit 0
